# systemd/anmin.service
[Unit]
Description = anmin

[Service]
ExecStart = /usr/local/sbin/anmin -q
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
